# concerto_full/DEBIAN/postinst
#!/bin/sh
set -e # lintian wants this script to exit when there are errors

#Source debconf library.
. /usr/share/debconf/confmodule

db_get concerto/passenger
passenger="$RET"

if [ "$passenger"  = "true" ]; then
  gem1.9.1 install --no-ri --no-rdoc passenger -v 4.0.5
  cd /var/lib/gems/1.9.1/gems/passenger-4.0.5
  rake apache2
fi

gem1.9.1 install --no-ri --no-rdoc rake bundler highline

#We'll leave database.yml alone if it exists already
if [ ! -f /usr/share/concerto/config/database.yml ];
then
  cp /usr/share/concerto/config/database.yml.mysql /usr/share/concerto/config/database.yml
fi
mkdir /usr/share/concerto/tmp
chmod -R 777 /usr/share/concerto/tmp
mkdir -p /usr/share/concerto/log 
chown -R www-data:www-data /usr/share/concerto
chmod 775 /usr/share/concerto/config/database.yml
chmod g+w /usr/share/concerto/log
chmod g+w /usr/share/concerto/public
#remove any old assets sitting around - Concerto will recompile
rm -rf /usr/share/concerto/public/assets

cd /usr/share/concerto
bundle install --path=vendor/bundle

a2ensite concerto

update-rc.d concerto defaults
invoke-rc.d concerto start

apache2ctl restart
